# deploy/systemd/openwebpanel-child@.service
[Unit]
Description=OpenWebPanel Child Daemon (%i)
Documentation=https://github.com/openwebcpanel/openwebcpanel
After=network.target

[Service]
Type=simple
User=openwebpanel
WorkingDirectory=/opt/openwebpanel/app
Environment=OWP_HOME_DIR=/opt/openwebpanel/homes/%i
Environment=OWP_CHILD_LISTEN=:9001
ExecStart=/opt/openwebpanel/app/bin/childd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
