# deploy/systemd/openwebpanel-child@.service
[Unit]
Description=OpenWebPanel Child Daemon (%i)
After=network.target

[Service]
Type=simple
User=claudeuser
WorkingDirectory=/home/claudeuser/openwebcpanel
Environment=OWP_HOME_DIR=/home/claudeuser/openwebcpanel/homes/%i
Environment=OWP_CHILD_LISTEN=:9001
ExecStart=/home/claudeuser/openwebcpanel/bin/childd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
